# roles/spark/templates/etc/systemd/system/spark-master.service
[Unit]
Description=Spark master
After=network.target

[Service]
Type=forking
ExecStart={{ spark_home }}/sbin/start-master.sh
ExecStop={{ spark_home }}/sbin/stop-master.sh
User={{ spark_user }}
Group={{ spark_group }}
Environment=SPARK_HOME={{ spark_home }}
Environment=JAVA_HOME={{ java_home }}


[Install]
WantedBy=multi-user.target
